# roles/redis/files/redisrw.service
[Unit]
Description=Load some dummy data into Redis

[Service]
Type=oneshot
ExecStart=/usr/local/bin/redisrw.sh
